# mail-relay.service
[Unit]
Description=Mail Relay - Relay de formularios de contacto multi-cliente
Documentation=https://github.com/x0-jonatanfp/mail-relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__SERVICE_DIR__

ExecStart=/usr/bin/node __SERVICE_DIR__/dist/index.js

Restart=on-failure
RestartSec=5
StartLimitBurst=5

Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
